# templates/upstart/zuul-server.conf
# zuul-server Zuul Server component
#

description     "Zuul Server"

# Make sure we start before an interface receives traffic
start on (starting network-interface
          or starting network-manager
          or starting networking)

stop on runlevel [!023456]

exec su -l -c "{{ zuul_server_bin }} -d -c {{ zuul_conf }}" {{ zuul_user }}
